# examples/p2p-im/tlsd-im-reconnect.timer
[Unit]
Description=TLSd P2P IM reconnect timer

[Timer]
OnUnitInactiveSec=600

[Install]
WantedBy=timers.target
